# var/lib/dpkg/info/yum.list
/.
/usr
/usr/bin
/usr/bin/yum
/usr/sbin
/usr/sbin/yum-updatesd
/usr/share
/usr/share/yum-cli
/usr/share/yum-cli/callback.py
/usr/share/yum-cli/cli.py
/usr/share/yum-cli/output.py
/usr/share/yum-cli/shell.py
/usr/share/yum-cli/utils.py
/usr/share/yum-cli/yumcommands.py
/usr/share/yum-cli/yummain.py
/usr/share/yum-cli/yumupd.py
/usr/share/man
/usr/share/man/man5
/usr/share/man/man5/yum.conf.5.gz
/usr/share/man/man5/yum-updatesd.conf.5.gz
/usr/share/man/man8
/usr/share/man/man8/yum-updatesd.8.gz
/usr/share/man/man8/yum.8.gz
/usr/share/man/man8/yum-shell.8.gz
/usr/share/locale
/usr/share/locale/ca
/usr/share/locale/ca/LC_MESSAGES
/usr/share/locale/ca/LC_MESSAGES/yum.mo
/usr/share/locale/cs
/usr/share/locale/cs/LC_MESSAGES
/usr/share/locale/cs/LC_MESSAGES/yum.mo
/usr/share/locale/da
/usr/share/locale/da/LC_MESSAGES
/usr/share/locale/da/LC_MESSAGES/yum.mo
/usr/share/locale/de
/usr/share/locale/de/LC_MESSAGES
/usr/share/locale/de/LC_MESSAGES/yum.mo
/usr/share/locale/es
/usr/share/locale/es/LC_MESSAGES
/usr/share/locale/es/LC_MESSAGES/yum.mo
/usr/share/locale/fr
/usr/share/locale/fr/LC_MESSAGES
/usr/share/locale/fr/LC_MESSAGES/yum.mo
/usr/share/locale/it
/usr/share/locale/it/LC_MESSAGES
/usr/share/locale/it/LC_MESSAGES/yum.mo
/usr/share/locale/ja
/usr/share/locale/ja/LC_MESSAGES
/usr/share/locale/ja/LC_MESSAGES/yum.mo
/usr/share/locale/ms
/usr/share/locale/ms/LC_MESSAGES
/usr/share/locale/ms/LC_MESSAGES/yum.mo
/usr/share/locale/nb
/usr/share/locale/nb/LC_MESSAGES
/usr/share/locale/nb/LC_MESSAGES/yum.mo
/usr/share/locale/pl
/usr/share/locale/pl/LC_MESSAGES
/usr/share/locale/pl/LC_MESSAGES/yum.mo
/usr/share/locale/pt
/usr/share/locale/pt/LC_MESSAGES
/usr/share/locale/pt/LC_MESSAGES/yum.mo
/usr/share/locale/pt_BR
/usr/share/locale/pt_BR/LC_MESSAGES
/usr/share/locale/pt_BR/LC_MESSAGES/yum.mo
/usr/share/locale/ru
/usr/share/locale/ru/LC_MESSAGES
/usr/share/locale/ru/LC_MESSAGES/yum.mo
/usr/share/locale/sr
/usr/share/locale/sr/LC_MESSAGES
/usr/share/locale/sr/LC_MESSAGES/yum.mo
/usr/share/locale/sr@latin
/usr/share/locale/sr@latin/LC_MESSAGES
/usr/share/locale/sr@latin/LC_MESSAGES/yum.mo
/usr/share/locale/sv
/usr/share/locale/sv/LC_MESSAGES
/usr/share/locale/sv/LC_MESSAGES/yum.mo
/usr/share/locale/zh_CN
/usr/share/locale/zh_CN/LC_MESSAGES
/usr/share/locale/zh_CN/LC_MESSAGES/yum.mo
/usr/share/doc
/usr/share/doc/yum
/usr/share/doc/yum/changelog.gz
/usr/share/doc/yum/FAQ.gz
/usr/share/doc/yum/README.source
/usr/share/doc/yum/README.Debian
/usr/share/doc/yum/README
/usr/share/doc/yum/TODO
/usr/share/doc/yum/AUTHORS
/usr/share/doc/yum/copyright
/usr/share/doc/yum/changelog.Debian.gz
/usr/share/doc/yum/PLUGINS.gz
/usr/share/pyshared
/usr/share/pyshared/yum
/usr/share/pyshared/yum/depsolve.py
/usr/share/pyshared/yum/sqlutils.py
/usr/share/pyshared/yum/misc.py
/usr/share/pyshared/yum/transactioninfo.py
/usr/share/pyshared/yum/sqlitesack.py
/usr/share/pyshared/yum/plugins.py
/usr/share/pyshared/yum/pgpmsg.py
/usr/share/pyshared/yum/parser.py
/usr/share/pyshared/yum/repos.py
/usr/share/pyshared/yum/rpmtrans.py
/usr/share/pyshared/yum/mdparser.py
/usr/share/pyshared/yum/config.py
/usr/share/pyshared/yum/metalink.py
/usr/share/pyshared/yum/logginglevels.py
/usr/share/pyshared/yum/Errors.py
/usr/share/pyshared/yum/i18n.py
/usr/share/pyshared/yum/history.py
/usr/share/pyshared/yum/constants.py
/usr/share/pyshared/yum/callbacks.py
/usr/share/pyshared/yum/__init__.py
/usr/share/pyshared/yum/packageSack.py
/usr/share/pyshared/yum/update_md.py
/usr/share/pyshared/yum/rpmsack.py
/usr/share/pyshared/yum/comps.py
/usr/share/pyshared/yum/repoMDObject.py
/usr/share/pyshared/yum/packages.py
/usr/share/pyshared/yum/failover.py
/usr/share/pyshared/yum/yumRepo.py
/usr/share/pyshared/rpmUtils
/usr/share/pyshared/rpmUtils/transaction.py
/usr/share/pyshared/rpmUtils/oldUtils.py
/usr/share/pyshared/rpmUtils/updates.py
/usr/share/pyshared/rpmUtils/__init__.py
/usr/share/pyshared/rpmUtils/miscutils.py
/usr/share/pyshared/rpmUtils/arch.py
/usr/share/python-support
/usr/share/python-support/yum.public
/usr/share/python-support/yum.private
/etc
/etc/yum
/etc/yum/repos.d
/etc/yum/yum.conf
/etc/yum/version-groups.conf
/etc/yum/yum-updatesd.conf
/etc/init.d
/etc/init.d/yum-updatesd
/etc/dbus-1
/etc/dbus-1/system.d
/etc/dbus-1/system.d/yum-updatesd.conf
/etc/logrotate.d
/etc/logrotate.d/yum
/etc/default
/etc/default/yum-updatesd
/var
/var/cache
/var/cache/yum
/var/lib
/var/lib/yum
